# deploy/com.heyim.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heyim.backend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mac/HeyIm/backend/.build/release/HeyImServer</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/mac/HeyIm/backend</string>
    
    <key>UserName</key>
    <string>mac</string>
    
    <key>GroupName</key>
    <string>staff</string>
    
    <key>StandardOutPath</key>
    <string>/Users/mac/HeyIm/logs/backend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mac/HeyIm/logs/backend-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>MODEL_PATH</key>
        <string>/Users/mac/HeyIm/models/RealisticVision_v51_split-einsum</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>Nice</key>
    <integer>-10</integer>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
