# deploy/com.heyim.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heyim.frontend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npm</string>
        <string>run</string>
        <string>start</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/mac/HeyIm/frontend</string>
    
    <key>UserName</key>
    <string>mac</string>
    
    <key>GroupName</key>
    <string>staff</string>
    
    <key>StandardOutPath</key>
    <string>/Users/mac/HeyIm/logs/frontend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mac/HeyIm/logs/frontend-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>5859</string>
        <key>NEXT_PUBLIC_API_URL</key>
        <string>http://localhost:5858</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
